# systemd/scheduled-reports.service
[Unit]
Description=Wolfe Scheduled Reports
After=network.target

[Service]
Type=simple
User=theplummer92
Environment=GOD_MODE_DATA_DIR=/home/theplummer92
Environment=GOD_MODE_ENV_FILE=/home/theplummer92/.env
WorkingDirectory=/home/theplummer92/god-mode-sync
ExecStart=/usr/bin/env TZ=America/New_York /home/theplummer92/venv/bin/python3 /home/theplummer92/god-mode-sync/scheduled_reports.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
