# systemd/sniper.service
[Unit]
Description=Wolfe Sniper Bot
After=network.target

[Service]
Type=simple
User=theplummer92
Environment=GOD_MODE_DATA_DIR=/home/theplummer92
Environment=GOD_MODE_ENV_FILE=/home/theplummer92/.env
WorkingDirectory=/home/theplummer92/god-mode-sync
ExecStart=/home/theplummer92/venv/bin/python3 /home/theplummer92/god-mode-sync/sniper_bot.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
